# daily-checkin/daily-checkin-block.service
[Unit]
Description=Run daily checkin block

[Service]
Type=oneshot
ExecStart=/home/pays0n/daily-checkin/block.sh
User=root
